# systemd/rfidreader.service
[Unit]
Description=PN532 RFID Reader Service (Waveshare HAT → UDP to Volumio)
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/data/rfidreader
ExecStart=/usr/bin/python3 /data/rfidreader/rfidreader.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
